# infra/onleaked-worker.conf
; /etc/supervisor/conf.d/onleaked-worker.conf
; Supervisor config for Onleaked OSINT queue workers

[program:onleaked-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/Projets/nealix/onleaked/artisan queue:work redis --queue=osint,default --tries=2 --timeout=240 --sleep=3 --max-time=3600
directory=/var/www/Projets/nealix/onleaked
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/Projets/nealix/onleaked/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
stopwaitsecs=300
